# systemd/user/battery-notify.service
[Unit]
Description=Battery level checker using Rofi

[Service]
Type=oneshot
ExecStart=%h/.config/battery-notify/check.sh

